# leaderboard/deploy/discord-bots-leaderboard.service
[Unit]
Description=Discord leaderboard bot service
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
WorkingDirectory=/home/ubuntu/discord-bots/leaderboard
EnvironmentFile=/home/ubuntu/leaderboard-secrets/app.env
Restart=on-failure
RestartSec=15s
ExecStart=/home/ubuntu/leaderboard-env/bin/python -m leaderboard.cli discord run
SyslogIdentifier=discord-bots-leaderboard

[Install]
WantedBy=multi-user.target